# build_scripts/db_setup/crontab
# Backup dspace database everyday at midnight
0 0 * * * sudo -u postgres -E sh /vagrant/build_scripts/db_setup/backup.sh